# cors-anywhere.service
[Unit]
Description=CORS Anywhere Proxy Server
After=network.target

[Service]
ExecStart=node /opt/cors-anywhere/server.js

[Install]
WantedBy=multi-user.target